# roles/interactive_tools/files/profile.d/sh.local
# for interactive shells only
if [[ $- == *i* ]]; then
  GREEN="\[$(tput setaf 2)\]"
  ORANGE="\[$(tput setaf 11)\]"
  RESET="\[$(tput sgr0)\]"
  PS1="${GREEN}[${ORANGE}\u@\h${RESET} \W${GREEN}]\$${RESET} "
  if command -v fzf > /dev/null; then
    eval "$(fzf --bash)"
  fi
  unset GREEN ORANGE RESET
fi
