# roles/grafana/files/grafana.container
[Unit]
Description=Grafana
After=victoriametrics.service

[Container]
Image=docker.io/grafana/grafana:latest
ContainerName=grafana
# Exposed for nginx proxy
PublishPort=3000:3000
Volume=grafana_data.volume:/var/lib/grafana
Volume=%h/.config/grafana/grafana.ini:/etc/grafana/grafana.ini:ro
Volume=%h/.config/grafana/provisioning:/etc/grafana/provisioning:ro

[Service]
Restart=always

[Install]
WantedBy=default.target
